# clean_logs.service
[Unit]
Description=Nettoyage des journaux système

[Service]
ExecStart=/usr/local/bin/TP1810/clean_logs.sh
